# system/systemupdate.service
[Unit]
Description=OTA

[Service]
ExecStart=/usr/bin/systemupdate

[Install]
WantedBy=graphical.target
